# Library/PrivateFrameworks/WebContentAnalysis.framework/Resources/com.apple.webcontentfilter.dns.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.webcontentfilter.dns</key>
		<true/>
	</dict>
	<key>Label</key>
	<string>com.apple.webcontentfilter.dns</string>
	<key>UserName</key>
	<string>nobody</string>
	<key>ProgramArguments</key>
	<array>
		<string>sandbox-exec</string>
		<string>-n</string>
		<string>/System/Library/PrivateFrameworks/WebContentAnalysis.framework/Resources/webfilterDNSd.sb</string>
		<string>/System/Library/PrivateFrameworks/WebContentAnalysis.framework/Resources/webfilterDNSd</string>
	</array>
	<key>EnableTransactions</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
</dict>
</plist>
