# Library/LaunchDaemons/com.apple.corestorage.corestoragehelperd.installer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.apple.corestorage.corestoragehelperd</string>
	<key>Program</key>
	<string>/usr/libexec/corestoragehelperd</string>
	<key>ProgramArguments</key>
	<array>
		<string>corestoragehelperd</string>
		<string>-i</string>
	</array>
	<key>EnableTransactions</key>
	<true/>
	<key>MachServices</key>
	<dict>
		<key>com.apple.corestorage.corestoragehelperd</key>
		<true/>
	</dict>
</dict>
</plist>
