# Library/LaunchDaemons/com.apple.wifi.WiFiAgent.BaseSystem.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.apple.wifi.WiFiAgent</string>
	<key>KeepAlive</key>
	<dict>
		<key>PathState</key>
		<dict>
			<key>/private/var/run/wifi-base-system</key>
			<true/>
		</dict>
		<key>Crashed</key>
		<true/>
	</dict>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/CoreServices/WiFiAgent.app/Contents/MacOS/WiFiAgent</string>
		<string>-BaseSystem</string>
	</array>
	<key>MachServices</key>
	<dict>
		<key>com.apple.wifi.WiFiAgent</key>
		<true/>
	</dict>
	<key>EnableTransactions</key>
	<true/>
	<key>ProcessType</key>
	<string>Adaptive</string>
</dict>
</plist>
