# Library/LaunchAgents/com.apple.TextInputMenuAgent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Program</key>
	<string>/System/Library/CoreServices/TextInputMenuAgent.app/Contents/MacOS/TextInputMenuAgent</string>
	<key>Label</key>
	<string>com.apple.TextInputMenuAgent</string>
	<key>RunAtLoad</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>KeepAlive</key>
	<dict>
		<key>Crashed</key>
		<true/>
	</dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>MachServices</key>
	<dict>
		<key>com.apple.tsm.uiserver</key>
		<dict>
			<key>HideUntilCheckIn</key>
			<true/>
		</dict>
	</dict>
	<key>LimitLoadToSessionType</key>
	<array>
		<string>Aqua</string>
	</array>
	<key>ProcessType</key>
	<string>App</string>
</dict>
</plist>
